# foundationdb/foundationdb.sv.conf
[program:fdbserver]
command=/bin/bash /usr/lib/foundationdb/foundationdb-start.sh
exitcodes=0
autorestart=false
startretries=0
startsecs=0
